# apps/gnosi/pipeline/backup_agents/com.gnosi.proactive-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosi.proactive-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/ismaelgarciafernandez/Projectes/monorepo/apps/gnosi/pipeline/skills/auto_improver/scripts/proactive_agent.py</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/ismaelgarciafernandez/Projectes/monorepo/apps/gnosi/pipeline/.tmp/proactive_agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ismaelgarciafernandez/Projectes/monorepo/apps/gnosi/pipeline/.tmp/proactive_agent.log</string>
</dict>
</plist>
